# TBD/init.service
[Unit]
Description=Initialize Display

[Service]
ExecStart=/bin/bash /opt/FiberPTS/init.sh

[Install]
WantedBy=multi-user.target