# TBD/screen.service
[Unit]
Description=Monitor and restart screen.py

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/src/screen/screen.py
User=root
Group=root
WorkingDirectory=/opt/FiberPTS
Environment="PATH=/usr/local/bin:/usr/bin"
PIDFile=/var/run/screen.pid
ExecStop=/bin/kill -15 \$MAINPID
TimeoutSec=5
Restart=always

[Install]
WantedBy=multi-user.target